# env-scripts/monthly.cron
# m h  dom mon dow   command
# At 1 am run the task once every day between the 5th and 10th of each month.
00 01 5-10 * * /bin/bash /usr/local/exec_monthly.sh >> /usr/local/data/exec_monthly.log 2>&1
